# default/replication.cron
0 3 * * * /usr/local/bin/replication.sh
